# supervisord.conf
[supervisord]
nodaemon=true

[program:prepare_project]
command=bash /prepare_project.sh
autostart=true

[program:dummy]
command=tail -f /dev/null
autostart=true
